# .config/systemd/user/hourly_background.timer
[Unit]
Description=hourlybackground

[Timer]
OnCalendar=*-*-* *:59:21
Persistent=true

[Install]
WantedBy=timers.target
