# .config/systemd/user/hourly_wallpaper_change.service
[Unit]
Description=Change wallpaper every hour

[Service]
Type=oneshot
ExecStart=/home/benkenobi/.local/bin/change_wallpaper
