# .config/systemd/user/hourly_background.service
[Unit]
Description=hourlybackground

[Service]
Type=simple
ExecStart=/home/benkenobi/.local/bin/change_background.py

